# bacula/platforms/solaris/bacula-dir.in
#! /bin/sh
#
# bacula       This shell script takes care of starting and stopping
#	       the bacula Director daemon
#
# chkconfig: 2345 92 9
# description: The Leading Open Source Backup Solution.
#
#  For Bacula release @VERSION@ (@DATE@) -- @DISTNAME@
#

DIR_USER=@dir_user@
DIR_GROUP=@dir_group@
DIR_OPTIONS=''

RETVAL=0
case "$1" in
    start)
       if [ "${DIR_USER}" != '' ]; then
          DIR_OPTIONS="${DIR_OPTIONS} -u ${DIR_USER}"
       fi
       if [ "${DIR_GROUP}" != '' ]; then
          DIR_OPTIONS="${DIR_OPTIONS} -g ${DIR_GROUP}"
       fi

       echo "Starting the Bacula Director: "
       @sbindir@/bacula-dir $2 ${DIR_OPTIONS} -c @sysconfdir@/bacula-dir.conf
       ;;
    stop)
       echo "Stopping the Director daemon: "
       pkill -x bacula-dir
       ;;
    restart)
       $0 stop
       sleep 5
       $0 start
       ;;
    *)
       echo "Usage: $0 {start|stop|restart}"
       exit 1
       ;;
esac
exit 0
