# bacula/platforms/systemd/bacula-sd.service.in
# Systemd Bacula Enterprise service file
#
# Copyright (C) Kern Sibbald 2011-2015
# License: BSD 2-Clause
#
# /lib/systemd/system/bacula-sd.service
#
# Description:
#    Used to start the bacula storage daemon service (bacula-sd)
#    enable : systemctl enable bacula-sd
#    start : systemctl start bacula-sd
#
#

# from http://www.freedesktop.org/software/systemd/man/systemd.unit.html
[Unit]
Description=Bacula Enterprise Storage Daemon service
Requires=network.target
After=network.target
RequiresMountsFor=@working_dir@ @sysconfdir@ @sbindir@ @piddir@

# from http://www.freedesktop.org/software/systemd/man/systemd.service.html
[Service]
Type=forking
User=@sd_user@
Group=@sd_group@
ExecStart=@sbindir@/bacula-sd -c @sysconfdir@/bacula-sd.conf
PIDFile=@piddir@/bacula-sd.@sd_port@.pid
StandardError=syslog

[Install]
WantedBy=multi-user.target
