# bacula/platforms/systemd/bacula-dir.service.in
# Systemd Bacula Enterprise service file
#
# Copyright (C) 2000-2015 by Kern Sibbald
# License: BSD 2-Clause; see file LICENSE-FOSS
#
# /lib/systemd/system/bacula-dir.service
#
# Description:
#    Used to start/stop/reload the bacula director daemon service (bacula-dir)
#
#    enable : systemctl enable bacula-dir
#    start : systemctl start bacula-dir
#
#

# From http://www.freedesktop.org/software/systemd/man/systemd.unit.html
[Unit]
Description=Bacula Enterprise Director Daemon service
Requires=network.target
After=network.target
RequiresMountsFor=@working_dir@ @sysconfdir@ @sbindir@ @piddir@

# From http://www.freedesktop.org/software/systemd/man/systemd.service.html
[Service]
Type=forking
User=@dir_user@
Group=@dir_group@
ExecStart=@sbindir@/bacula-dir -c @sysconfdir@/bacula-dir.conf
PIDFile=@piddir@/bacula-dir.@dir_port@.pid
ExecReload=@sbindir@/bacula-dir -t -c @sysconfdir@/bacula-dir.conf
ExecReload=/bin/kill -HUP $MAINPID
StandardError=syslog

[Install]
WantedBy=multi-user.target
